# srcs/nginx/srcs/supervisord.conf
[supervisord]
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=false
startretries=3

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
startretries=3

[include]
files = /etc/supervisor/conf.d/*.conf
